# luci-app-tinyfilemanager/root/etc/uci-defaults/tinyfilemanager.sh
#!/bin/sh

cd /
ln -s / /www/tinyfilemanager/rootfs
exit 0
