# luci-app-tinyfilemanager/root/etc/uci-defaults/tinyfilemanager.sh
#!/bin/sh

cd /
[ ! -d /www/tinyfilemanager/rootfs ] && ln -s / /www/tinyfilemanager/rootfs
exit 0
